# apps/runtime-xpra/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/sandbox/supervisord.log
pidfile=/run/sandbox/supervisord.pid

[program:xpra]
command=/bin/bash /opt/sandbox/scripts/start_xpra.sh
priority=5
autorestart=true
stdout_logfile=/var/log/sandbox/xpra.log
stderr_logfile=/var/log/sandbox/xpra.err.log

[program:fcitx]
command=/bin/bash /opt/sandbox/scripts/start_fcitx.sh
priority=10
autorestart=true
stdout_logfile=/var/log/sandbox/fcitx.log
stderr_logfile=/var/log/sandbox/fcitx.err.log
environment=DISPLAY=":100"

[program:chromium]
command=/bin/bash /opt/sandbox/scripts/start_browser.sh
priority=20
autorestart=true
stdout_logfile=/var/log/sandbox/chromium.log
stderr_logfile=/var/log/sandbox/chromium.err.log
environment=DISPLAY=":100"

[program:cdp-proxy]
command=/bin/bash /opt/sandbox/scripts/start_cdp_proxy.sh
priority=25
autorestart=true
stdout_logfile=/var/log/sandbox/cdp-proxy.log
stderr_logfile=/var/log/sandbox/cdp-proxy.err.log
