# apps/sandbox-runtime/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/sandbox/supervisord.log
pidfile=/run/sandbox/supervisord.pid

[program:xvfb]
command=/bin/bash /opt/sandbox/scripts/start_x.sh
priority=5
autorestart=true
stdout_logfile=/var/log/sandbox/xvfb.log
stderr_logfile=/var/log/sandbox/xvfb.err.log

[program:fcitx]
command=/bin/bash /opt/sandbox/scripts/start_fcitx.sh
priority=10
autorestart=true
stdout_logfile=/var/log/sandbox/fcitx.log
stderr_logfile=/var/log/sandbox/fcitx.err.log
environment=DISPLAY=":99"

[program:chromium]
command=/bin/bash /opt/sandbox/scripts/start_browser.sh
priority=20
autorestart=true
stdout_logfile=/var/log/sandbox/chromium.log
stderr_logfile=/var/log/sandbox/chromium.err.log
environment=DISPLAY=":99"

[program:cdp-proxy]
command=/bin/bash /opt/sandbox/scripts/start_cdp_proxy.sh
priority=25
autorestart=true
stdout_logfile=/var/log/sandbox/cdp-proxy.log
stderr_logfile=/var/log/sandbox/cdp-proxy.err.log

[program:x11vnc]
command=/bin/bash /opt/sandbox/scripts/start_x11vnc.sh
priority=30
autorestart=true
stdout_logfile=/var/log/sandbox/x11vnc.log
stderr_logfile=/var/log/sandbox/x11vnc.err.log
environment=DISPLAY=":99"

[program:websockify]
command=/bin/bash /opt/sandbox/scripts/start_websockify.sh
priority=40
autorestart=true
stdout_logfile=/var/log/sandbox/websockify.log
stderr_logfile=/var/log/sandbox/websockify.err.log
environment=DISPLAY=":99"
